# etc/init.d/heizungpanel
#!/bin/sh /etc/rc.common
USE_PROCD=1
START=95
STOP=10
BUILD_TAG="commit:8b755f2"

start_service() {
  . /lib/functions.sh
  config_load heizungpanel

  local can_if can_bitrate can_setup listen_only write_mode
  local mqtt_host mqtt_port mqtt_base
  local publish_raw publish_state

  config_get can_if          main can_if          "can0"
  config_get can_bitrate     main can_bitrate     "69144"
  config_get can_setup       main can_setup       "1"
  config_get write_mode      main write_mode      "0"
  config_get publish_raw     main publish_raw     "1"
  config_get publish_state   main publish_state   "0"
  [ "$write_mode" = "1" ] && listen_only="0" || listen_only="1"

  config_get mqtt_host   main mqtt_host   "127.0.0.1"
  config_get mqtt_port   main mqtt_port   "1883"
  config_get mqtt_base   main mqtt_base   "heizungpanel"

  local topic_raw="${mqtt_base}/raw"
  local topic_state="${mqtt_base}/state"

  logger -t heizungpanel "service start ($BUILD_TAG)"

  if [ "$can_setup" = "1" ]; then
    case "$can_if" in
      can*|vcan*|slcan*) ;;
      *)
        logger -t heizungpanel "Refusing CAN setup on non-CAN interface: $can_if (set option can_setup '0' or can_if 'can0')"
        return 1
        ;;
    esac

    local lo_arg="listen-only off"
    [ "$listen_only" = "1" ] && lo_arg="listen-only on"

    if ! ip link show "$can_if" >/dev/null 2>&1; then
      logger -t heizungpanel "CAN interface not found: $can_if"
      return 1
    fi

    ip link set "$can_if" down >/dev/null 2>&1 || true
    if ! ip link set "$can_if" type can bitrate "$can_bitrate" $lo_arg >/dev/null 2>&1; then
      logger -t heizungpanel "CAN setup failed: $can_if bitrate $can_bitrate $lo_arg"
      return 1
    fi
    if ! ip link set "$can_if" up >/dev/null 2>&1; then
      logger -t heizungpanel "CAN bring-up failed: $can_if"
      return 1
    fi
    logger -t heizungpanel "CAN setup ok: $can_if bitrate $can_bitrate ($lo_arg)"
  fi

  if [ "$publish_raw" = "1" ]; then
    procd_open_instance raw
    procd_set_param command /usr/libexec/heizungpanel/raw_bridge.sh "$can_if" "$mqtt_host" "$mqtt_port" "$topic_raw"
    procd_set_param respawn 3600 5 5
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_close_instance
  fi

  if [ "$publish_state" = "1" ]; then
    procd_open_instance state
    procd_set_param command /usr/libexec/heizungpanel/state_bridge.sh "$can_if" "$can_bitrate" "$mqtt_host" "$mqtt_port" "$topic_state"
    procd_set_param respawn 3600 5 5
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_close_instance
  fi

  local summary="MQTT publishing: raw=${publish_raw}:${topic_raw} bootstrap(file)=/tmp/heizungpanel/bootstrap.json"
  summary="$summary state(legacy,optional)=${publish_state}:${topic_state}"
  logger -t heizungpanel "$summary"
}
